# dev1/services/chatterbox-dev1.service
[Unit]
Description=Chatterbox Application - Dev1 Machine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/chatterbox/dev1
ExecStart=/usr/bin/python3 /opt/chatterbox/dev1/client/main.py --config /opt/chatterbox/dev1/config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/chatterbox/dev1/logs
ReadOnlyPaths=/var/log

MemoryMax=512M
CPUQuota=50%

KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
